# .config/systemd/user/taskd_client.service
[Unit]
Description=Runs a synchronization of taskwarrior

[Service]
ExecStart=/usr/bin/task sync

